# shader-overlay/flick-effects.service
[Unit]
Description=Flick Effects Overlay
After=graphical-session.target
PartOf=graphical-session.target

[Service]
Type=simple
ExecStart=%h/.local/bin/flick-effects
Restart=on-failure
RestartSec=3

[Install]
WantedBy=graphical-session.target
